# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
loglevel=info

[program:xvfb]
command=Xvfb :99 -screen 0 %(ENV_RESOLUTION)s
priority=10
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:fluxbox]
command=fluxbox -display :99
priority=15
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/bin/bash -c 'x11vnc -display :99 -rfbport %(ENV_VNC_PORT)s -forever -shared $([ -n "$VNC_PASSWORD" ] && echo "-passwd $VNC_PASSWORD" || echo "-nopw")'
priority=20
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:novnc]
command=/bin/bash -c "websockify --web=/usr/share/novnc %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s"
priority=20
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:chrome]
command=/opt/bin/chrome-launcher.sh
priority=30
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
user=chrome
environment=HOME="/home/chrome",DISPLAY=":99"
